# .config/systemd/user/obsidian-sync.service
[Unit]
Description=Obsidian rclone sync

[Service]
Type=oneshot
ExecStart=/home/furuike/.config/obsidian-sync/sync.sh
